# .config/containers/systemd/crowdsec.container
[Unit]
Description=CrowdSec Security Engine
After=network-online.target stream-pod.service
Wants=stream-pod.service

[Container]
ContainerName=crowdsec
Image=docker.io/crowdsecurity/crowdsec:v1.7.4-slim
Pod=stream.pod

EnvironmentFile=/srv/quartermaster/plex/crowdsec.env

Volume=/srv/quartermaster/plex/data/crowdsec/config:/etc/crowdsec:Z
Volume=/srv/quartermaster/plex/data/crowdsec/data:/var/lib/crowdsec/data:Z
Volume=/srv/quartermaster/plex/data/logs:/var/log/caddy:ro,z
Volume=/srv/quartermaster/plex/config/acquis.yaml:/etc/crowdsec/acquis.yaml:ro,z

HealthCmd=cscli lapi status
HealthInterval=30s
HealthRetries=3
HealthStartPeriod=60s
HealthTimeout=10s

NoNewPrivileges=true

[Install]
WantedBy=default.target
